# desktop/wake-word/com.jarvis.morning.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  Jarvis Morning Briefing — launchd agent (Fase 5)
  Roda todo dia às 8h00.

  Instalar:
    cp com.jarvis.morning.plist ~/Library/LaunchAgents/
    launchctl load ~/Library/LaunchAgents/com.jarvis.morning.plist

  Testar agora (sem esperar 8h):
    launchctl start com.jarvis.morning

  Desinstalar:
    launchctl unload ~/Library/LaunchAgents/com.jarvis.morning.plist
    rm ~/Library/LaunchAgents/com.jarvis.morning.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.jarvis.morning</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/matheussilva/Documents/Developer2026/Jarvis/Jarvis/desktop/wake-word/.venv/bin/python3</string>
        <string>-u</string>
        <string>/Users/matheussilva/Documents/Developer2026/Jarvis/Jarvis/desktop/wake-word/jarvis_morning.py</string>
    </array>

    <!-- Roda todo dia às 8:00 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>8</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- KeepAlive FALSE — roda uma vez e para (Resource Safety) -->
    <key>KeepAlive</key>
    <false/>

    <key>StandardOutPath</key>
    <string>/tmp/jarvis-morning.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/jarvis-morning.log</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
